# .config/systemd/user/udiskie.service
[Unit]
Description=Udiskie Automounter Daemon

[Service]
ExecStart=/usr/bin/udiskie
Restart=always

[Install]
WantedBy=default.target
