# .local/etc/systemd/user/earlyoom.service
[Unit]
Description=Early Out Of Memory Killer

[Service]
ExecStart=/usr/bin/earlyoom

[Install]
WantedBy=default.target
